# services/init.d/breathgslb
#!/sbin/openrc-run

name="BreathGSLB"
description="Authoritative GSLB DNS"

# Customize these variables to match your environment
# User account the service runs as
SERVICE_USER="breathgslb"
# Group the service runs as
SERVICE_GROUP="breathgslb"
# Directory where log files are written
LOG_DIR="/var/log/breathgslb"
# Log file for both stdout and stderr
LOG_FILE="$LOG_DIR/breathgslb.log"
# Location of the PID file
PID_FILE="/run/breathgslb.pid"

# Path to the BreathGSLB executable
command="/usr/local/bin/breathgslb"
# Path to the configuration file
CONFIG_FILE="/etc/breathgslb/config.yaml"
# Base64‑encoded license payload (optional)
LICENSE_PAYLOAD=""
# Additional custom flags (optional)
EXTRA_ARGS=""
# Command line arguments
command_args="-config ${CONFIG_FILE}"
[ -n "${LICENSE_PAYLOAD}" ] && command_args="${command_args} -license-payload ${LICENSE_PAYLOAD}"
[ -n "${EXTRA_ARGS}" ] && command_args="${command_args} ${EXTRA_ARGS}"
# Run the command as the configured user and group
command_user="${SERVICE_USER}:${SERVICE_GROUP}"

# Use OpenRC's supervise-daemon for process supervision
supervisor="supervise-daemon"
# PID file used by supervise-daemon
pidfile="$PID_FILE"
# Log files used by supervise-daemon
output_log="$LOG_FILE"
error_log="$LOG_FILE"

# advertise custom commands
extra_started_commands="reload"
description_reload="Reload configuration (SIGHUP)"

depend() {
        # Ensure the network stack is available before starting
        need net
}

start_pre() {
        # Create the log directory with appropriate ownership and permissions
        checkpath -d -o "${SERVICE_USER}:${SERVICE_GROUP}" -m 0755 "$LOG_DIR"
}

reload() {
        if ! service_started ; then
                eerror "$name is not running"
                return 1
        fi
        ebegin "Reloading $name (SIGHUP)"
        # Rely on supervise-daemon to send SIGHUP to the running process
        /sbin/supervise-daemon "$RC_SVCNAME" --signal HUP
        eend $?
}

