# services/init.d/greathgslb
#!/sbin/openrc-run
name="BreathGSLB"
description="Authoritative GSLB DNS"

command="/usr/local/bin/breathgslb"
command_args="-config /etc/breathgslb/config.yaml"
command_user="breathgslb:breathgslb"

supervisor="supervise-daemon"
pidfile="/run/breathgslb.pid"

output_log="/var/log/breathgslb/breathgslb.log"
error_log="/var/log/breathgslb/breathgslb.log"

# advertise custom commands
extra_started_commands="reload"
description_reload="Reload configuration (SIGHUP)"

depend() {
	need net
}

start_pre() {
	checkpath -d -o breathgslb:breathgslb -m 0755 /var/log/breathgslb
}

reload() {
	if ! service_started ; then
		eerror "$name is not running"
		return 1
	fi
	ebegin "Reloading $name (SIGHUP)"
	# prefer supervise-daemon signal; fall back to start-stop-daemon
	/sbin/supervise-daemon "$RC_SVCNAME" --signal HUP \
	|| start-stop-daemon --signal HUP --pidfile "$pidfile"
	eend $?
}

