# services/systemd/breathgslb.service
[Unit]
# Brief description of the service.
Description=Breath GSLB Authoritative DNS
# Start only after networking is fully configured.
After=network-online.target
# Ensure the network stack is pulled in as a dependency.
Wants=network-online.target

[Service]
# Load optional environment overrides.
EnvironmentFile=/etc/breathgslb/env
# Default path to the configuration file; can be overridden in env file.
Environment=BREATHGSLB_CONFIG=/etc/breathgslb/config.yaml
# Launch the daemon with the configured file and any extra flags.
ExecStart=/usr/local/bin/breathgslb -config ${BREATHGSLB_CONFIG} ${BREATHGSLB_FLAGS}
# Allow reloading configuration without full restart.
ExecReload=/bin/kill -HUP $MAINPID
# Run under a dedicated unprivileged user.
User=breathgslb
# Dedicated group for file and process ownership.
Group=breathgslb
# Permit binding to privileged ports such as 53.
AmbientCapabilities=CAP_NET_BIND_SERVICE
# Limit available capabilities to reduce attack surface.
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# Disallow gaining additional privileges after start.
NoNewPrivileges=true
# Restart automatically if the process exits unexpectedly.
Restart=on-failure

[Install]
# Install target that starts the service during normal boot.
WantedBy=multi-user.target

